# etc/supervisor/conf.d/celery.conf
[program:celery]

# Set Django environment variables
environment=DJANGO_SETTINGS_MODULE="rosreestr.settings"

# Directory where your site's project files are located
directory=/home/web/app

command=/home/web/app/venv/bin/celery -A rosreestr worker --loglevel="INFO"

# Number of processes to startup, roughly the number of CPUs you have
numprocs=1

# Give each process a unique name so they can be told apart
process_name=celery%(process_num)d

# Automatically start and recover processes
autostart=true
autorestart=true

# Choose where you want your log to go
stdout_logfile=/var/log/celery.log
redirect_stderr=true
